# roles/13-gnodeb/templates/enodeb.service.j2
[Unit]
Description=Enodeb Service

[Service]
ExecStart=/usr/local/bin/enodeb.sh

[Install]
WantedBy=default.target