# roles/14-ue/templates/ue.service.j2
[Unit]
Description=UE Service

[Service]
ExecStart=/usr/local/bin/ue.sh

[Install]
WantedBy=default.target